# ansible/roles/infrastructure/templates/dashboard-sync.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Periodic reconcile of Traefik routes + SSO provider config

[Timer]
OnBootSec=2min
OnCalendar={{ dashboard_sync_interval }}
RandomizedDelaySec=60s
Persistent=true
Unit=dashboard-sync.service

[Install]
WantedBy=timers.target
