# ansible/roles/backup/templates/catena-backup.service.j2
[Unit]
Description=catena restic backup (manual oneshot -- no timer in CE)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
# systemd services run without HOME by default. Restic needs it (or
# XDG_CACHE_HOME) to keep a local metadata cache -- without a cache
# it re-downloads repo indexes on every run, which is slow and burns
# S3 egress. Point both at a persistent state dir so the cache
# survives reboots.
Environment=HOME=/var/lib/catena/restic-home
Environment=XDG_CACHE_HOME=/var/lib/catena/restic-cache
EnvironmentFile={{ backup_env_file }}
ExecStartPre=/bin/mkdir -p /var/lib/catena/restic-home /var/lib/catena/restic-cache
# R21: mutex with auto-update / stack-update via flock on /run/catena.lock.
# flock holds the lock for the LIFETIME of the wrapped command and releases
# it on exit OR crash. The prior ExecStartPre PID-file helper recorded the
# pre-step's own PID, which died the instant ExecStartPre returned -- so a
# contender saw a "stale" lock and acquired immediately, and nothing ever
# released it. -w 7200: wait up to 2h for a holder, then flock exits non-zero
# (systemd failure -> dead-man alert) rather than run concurrently.
ExecStart=/usr/bin/flock -w 7200 /run/catena.lock {{ backup_wrapper_script }}
# Kill the run if it exceeds 4 hours -- better to fail loudly (and ping /fail)
# than to silently run for days.
TimeoutStartSec=4h
# Nice the CPU + IO so backups don't starve foreground workloads.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
